# collection/roles/cluster_node/templates/swapoff.service
[Unit]

[Service]
Type=oneshot
RemainAfterExit=yes
ExecStart=/sbin/swapoff -a

[Install]
WantedBy=default.target
