# collection/roles/service_keepalived/templates/load_ip_vs.service
[Unit]

[Service]
Type=oneshot
RemainAfterExit=yes
ExecStart=/sbin/modprobe -a ip_vs

[Install]
WantedBy=default.target
