# roles/service_docker/templates/docker-prune.timer
[Unit]
Description=Prune useless Docker resources every week

[Timer]
OnCalendar=weekly

[Install]
WantedBy=timers.target
