# collection/roles/service_docker/templates/docker-prune.service
[Unit]
Description=Prune useless Docker resouces

[Service]
Type=simple
ExecStart=/usr/local/bin/docker-prune
